# packer-scripts/ansible/roles/fluentd/templates/supervisor_fluentd.conf.j2
[program:fluentd]
command=/usr/sbin/td-agent --no-supervisor
